# conf/runit/nginx/run
#!/usr/bin/env bash

exec /usr/sbin/nginx -g "daemon off;"
